# retrofit-capture.service
[Unit]
Description=RetroFit Image Capture Service v2.1
# Repository: https://github.com/aditya08deole/Flow-1.git
After=network-online.target
Wants=network-online.target
StartLimitBurst=3
StartLimitIntervalSec=300

[Service]
Type=simple
User=__SERVICE_USER__
WorkingDirectory=__INSTALL_DIR__

# Python environment
Environment="PYTHONDONTWRITEBYTECODE=1"
Environment="PYTHONUNBUFFERED=1"
Environment="VIRTUAL_ENV=__INSTALL_DIR__/.venv"
Environment="PATH=__INSTALL_DIR__/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

# Pre-start: wait for network stack to fully settle
ExecStartPre=/bin/sleep 10

# Start using .venv Python
ExecStart=__INSTALL_DIR__/.venv/bin/python3 __INSTALL_DIR__/main_service.py

# Restart policy: only restart on non-fatal failures, stop after 3 rapid crashes
Restart=on-failure
RestartSec=30
RestartPreventExitStatus=2 3

# Resource limits (Pi Zero W safety)
MemoryMax=300M
CPUQuota=60%

[Install]
WantedBy=multi-user.target
